# piscripts/51-android.rules
SUBSYSTEM=="usb", ATTR{idVendor}=="04e8", RUN+="/usr/local/bin/backup.sh", OWNER="pi" 
